# tone-10hz/usr/local/lib/systemd/user/tone-10hz.timer
[Unit]
Description=Play tone-10hz.wav every 25 minutes

[Timer]
OnBootSec=25min
OnUnitActiveSec=25min

[Install]
WantedBy=default.target
